# templates/ssh-socks.socket
[Unit]
Description=Socket for SOCKS proxy to %host%

[Socket]
ListenStream=127.0.0.1:%port%

[Install]
WantedBy=sockets.target
