# templates/ssh-socks-proxy.service
[Unit]
Description=SOCKS proxy to %host%

[Service]
Type=forking
ExecStart=/usr/bin/ssh -o ExitOnForwardFailure=yes -f -NTCD 127.0.0.1:%port% %host%
