# templates/ssh-socks-proxy.service
[Unit]
Description=SOCKS proxy to %host%

[Service]
ExecStart=/usr/bin/ssh -NTD 127.0.0.1:%port% %host%
